# files/lib/systemd/system/yafti.service
[Unit]
Description=Yafti First Run

[Service]
ExecStart=/usr/local/bin/firstrun.sh

[Install]
WantedBy=multi-user.target
