# files/usr/lib/systemd/user/wallpaper-cycler.timer
[Unit]
Description=Run change-wallapper periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
