# files/etc/xdg/autostart/onboarding.desktop
[Desktop Entry]
Type=Application
Name=Onboarding
Exec=onboarding
OnlyShowIn=GNOME;
AutostartCondition=unless-exists onboarding/disable